# 20190328 Ansible Part 01/templates/app.ini
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:app]
directory={{ ansible_user_dir }}
command={{ ansible_user_dir }}/.virtualenv/bin/gunicorn app:app -b 127.0.0.1:5000
autostart=true
autorestart=true
